# docker/php/supervisord.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=4
redirect_stderr=true
stopwaitsecs=3600
stdout_logfile=/app/storage/logs/queue.log

[program:laravel-s]
directory=/app
command=/usr/local/bin/php bin/laravels start -i
redirect_stderr=true
autostart=true
autorestart=true
numprocs=1
process_name=%(program_name)s_%(process_num)s
stdout_logfile=/app/storage/logs/laravel-s.log

